# etc/profile.d/FreeBASIC-1.07.3.sh
export PATH=$PATH:/usr/local/FreeBASIC-1.07.3/bin
