# etc/profile.d/geant4.sh
. /usr/local/geant4.10.03.p01/bin/geant4.sh
. /usr/local/geant4.10.03.p01/share/Geant4-10.3.1/geant4make/geant4make.sh
