# etc/profile.d/GEF-2021-V1-1.sh
export PATH=$PATH:/usr/local/GEF-2021-V1-1
