# etc/profile.d/gnu_octave.sh
export PATH=/usr/local/octave-5.2.0/bin:$PATH
