# etc/profile.d/serpent-2.1.32.sh
export PATH=/usr/local/serpent-2.1.32:$PATH
export SERPENT_DATA=/usr/local/serpent-2.1.32/xsdata
export SERPENT_ACELIB=sss_endfb7u.xsdata
